# deploy/systemd/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for Pollex API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared tunnel --config /home/manu/.cloudflared/config.yml run
Restart=on-failure
RestartSec=10

# Note: User=manu and hardening directives (ProtectSystem, ProtectHome) are intentionally
# omitted — they cause failures on JetPack 4.6 (Ubuntu 18.04) systemd. See ADR-008.

[Install]
WantedBy=multi-user.target
